# contrib/init/mulecoind.conf
description "Mulecoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env MULECOIND_BIN="/usr/bin/mulecoind"
env MULECOIND_USER="mulecoin"
env MULECOIND_GROUP="mulecoin"
env MULECOIND_PIDDIR="/var/run/mulecoind"
# upstart can't handle variables constructed with other variables
env MULECOIND_PIDFILE="/var/run/mulecoind/mulecoind.pid"
env MULECOIND_CONFIGFILE="/etc/mulecoin/mulecoin.conf"
env MULECOIND_DATADIR="/var/lib/mulecoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # mulecoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$MULECOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run mulecoind."
        echo "The setting must appear in $MULECOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Mulecoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$MULECOIND_PIDDIR"
    chmod 0755 "$MULECOIND_PIDDIR"
    chown $MULECOIND_USER:$MULECOIND_GROUP "$MULECOIND_PIDDIR"
    chown $MULECOIND_USER:$MULECOIND_GROUP "$MULECOIND_CONFIGFILE"
    chmod 0660 "$MULECOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$MULECOIND_PIDFILE" \
    --chuid $MULECOIND_USER:$MULECOIND_GROUP \
    --exec "$MULECOIND_BIN" \
    -- \
    -pid="$MULECOIND_PIDFILE" \
    -conf="$MULECOIND_CONFIGFILE" \
    -datadir="$MULECOIND_DATADIR" \
    -disablewallet \
    -daemon

